# autostartpythontestbot.service
[Unit]
Description=pythontestbot

[Service]
ExecStart=/root/pythontestbot/runpythontestbot.sh

[Install]
WantedBy=multi-user.target